# fireplace.service.in
[Unit]
Description=Fireplace is the B2G Frame Host
ConditionUser=!root

[Service]
Type=simple
Environment=HYBRIS_EGLPLATFORM=null
ExecStartPre=-/bin/sh -ec '/bin/echo QUIT >/run/psplash_fifo'
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/fireplace
Restart=always
